# recipes-graphics/hyperion/files/init/systemd/hyperionaml.timer
[Unit]
Description=Starts hyperionaml

[Timer]
OnBootSec=1min
Unit=hyperionaml.service

[Install]
WantedBy=multi-user.target
